# root/etc/uci-defaults/luci-fluxproxy-migration
#!/bin/sh

ucode "/etc/fluxproxy/scripts/migrate_config.uc"

exit 0
